# ToyVpnServer/env.service
[Unit]
Description=environment for ToyVpnServer
After=network.target

[Service]
ExecStart=/home/abhishek/Downloads/ToyVpnServer-master/env.sh

[Install]
WantedBy=multi-user.target





















































